# core/automation/launchd-plists/com.brain.email-triage-v2-evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brain.email-triage-v2-evening</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/homeserver/Obsidian/personal-os/core/automation/email-triage-v2.sh</string>
        <string>--run-type</string>
        <string>pm</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/homeserver/Obsidian/personal-os</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>15</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/Users/homeserver/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/homeserver</string>
        <key>USER</key>
        <string>homeserver</string>
        <key>TRIAGE_V2_PROVIDER_MODE</key>
        <string>gmail</string>
        <key>TRIAGE_V2_SENDER_MODE</key>
        <string>gmail</string>
        <key>TRIAGE_V2_DRAFT_MODE</key>
        <string>superhuman_preferred</string>
        <key>TRIAGE_V2_SUPERHUMAN_ENABLED</key>
        <string>1</string>
        <key>TRIAGE_V2_ALWAYS_RECONCILE</key>
        <string>1</string>
        <key>TRIAGE_V2_ACCOUNTS</key>
        <string>work,personal</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/homeserver/Obsidian/personal-os/logs/email-triage-v2-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/homeserver/Obsidian/personal-os/logs/email-triage-v2-launchd-error.log</string>
</dict>
</plist>
